# hackpsf.service
[Unit]
Description=Hack PSF

[Service]
Type=oneshoot
ExecStart=/usr/bin/hackpsf

[Install]
WantedBy=multi-user.target
